# resources/com/google/gerrit/pgm/init/gerrit.socket
[Unit]
Description=Gerrit HTTP socket

[Socket]
ListenStream=80
Accept=no

[Install]
WantedBy=sockets.target
